# etc/systemd/system/docker.service.d/custom-environment-file.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// $OPTIONS
EnvironmentFile=-/etc/sysconfig/docker
